# templates/supervisor-client.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startretries=1
numprocs=1
startsecs=0

[program:ancile]
directory=/ancile
command=/ancile/client/start.sh
autostart=true
numproces=1
startsecs=0
